# DebBuild/usenet-systemd/DEBIAN/postinst
invoke-rc.d usenet-inn2 start || exit $?

exit 0