# DebBuild/usenet-systemd/DEBIAN/postinst
systemctl enable usenet-inn2.service
systemctl start usenet-inn2.service